# usr/lib/systemd/system/buildrust.timer
[Unit]
Description=Run the rust builds nightly

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
