# usr/lib/systemd/system/buildrust.service
[Unit]
Description=BuildRust

[Service]
Type=simple
ExecStart=/usr/local/bin/BuildRust.sh
